# backups/v1.0_stable_20251214/scripts/crontab
0 8 * * * cd /app && /usr/local/bin/python scripts/daily_pipeline.py >> /var/log/cron.log 2>&1