# docker/supervisord.conf
# =================================================================
# Supervisor生产环境配置
# Production Supervisor Configuration
# =================================================================

[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:football-prediction]
command=gunicorn --config /etc/gunicorn.conf.py src.main:app
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=ENV="production",PYTHONPATH="/app"
stopwaitsecs=10
stopsignal=INT
killasgroup=true
priority=999

[program:worker]
command=celery -A src.core.celery worker --loglevel=info --concurrency=4
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=ENV="production",PYTHONPATH="/app"
stopwaitsecs=10
stopsignal=TERM
killasgroup=true
priority=998

[program:beat]
command=celery -A src.core.celery beat --loglevel=info
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/beat.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=ENV="production",PYTHONPATH="/app"
stopwaitsecs=10
stopsignal=TERM
killasgroup=true
priority=997

[group:football-prediction-group]
programs=football-prediction,worker,beat
priority=999